# Tools/SkyrimVR/systemd/skyrim-together-vr-build-cleanup.service
[Unit]
Description=Clean stale Skyrim Together VR build storage

[Service]
Type=oneshot
ExecStart=%h/Documents/SkyrimModding/SkyrimTogetherVR/Tools/SkyrimVR/cleanup_build_storage.sh --scheduled --max-age-days 2 --trim --rebuildable-caches
